# services/payload.service
[Unit]
Description=LinkVST headless host
After=sound.target

[Service]
Environment=JACK_NO_AUDIO_RESERVATION=1
User=pi
Group=pi
WorkingDirectory=/home/pi/vcpi
RuntimeDirectory=linkvst
ExecStartPre=/bin/sleep 10
ExecStart=/home/pi/vcpi/venv/bin/python vst_host.py serve --sock /run/linkvst/linkvst.sock
Restart=always
RestartSec=10s
KillMode=process
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target
